# checkupdates.timer
[Unit]
Description=Run checkupdates periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
